# root/etc/s6-overlay/s6-rc.d/init-miniflux-config/up
/etc/s6-overlay/s6-rc.d/init-miniflux-config/run
